# systemd/raspimidihub.service
[Unit]
Description=RaspiMIDIHub - Automatic USB MIDI Hub
After=sound.target network.target
Wants=sound.target

[Service]
Type=notify
ExecStart=/usr/bin/python3 -m raspimidihub
Restart=always
RestartSec=2
WatchdogSec=30
NotifyAccess=main

# Run as root for ALSA sequencer, LED, hostapd, dnsmasq
User=root

# Hardening (relaxed for hostapd/dnsmasq child processes)
ProtectHome=true
PrivateTmp=false
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=false

[Install]
WantedBy=multi-user.target
