# systemd/raspimidihub-hostapd.service
[Unit]
Description=hostapd for RaspiMIDIHub access point
Documentation=man:hostapd(8)
After=network.target
# Don't try to start without our generated config in place. raspimidihub
# writes /run/raspimidihub/hostapd.conf before invoking systemctl start;
# this guards against an empty config blocking the AP at first boot.
ConditionPathExists=/run/raspimidihub/hostapd.conf

[Service]
# Type=simple keeps hostapd in foreground so its stderr lands in
# journald automatically — no logger_syslog workaround needed.
Type=simple
ExecStart=/usr/sbin/hostapd /run/raspimidihub/hostapd.conf
Restart=on-failure
RestartSec=2
# Keep hostapd on the housekeeping cores (0-1): core 3 is reserved
# for the asyncio main loop and core 2 for the plugin pool. hostapd
# needs a core with normal timer ticks anyway (it can't beacon on a
# nohz_full isolated core), and 0-1 have headroom for the WiFi driver
# IRQs, hostapd, dnsmasq, and the rest of userspace.
CPUAffinity=0-1

# No [Install] section: only raspimidihub.service starts us via
# `systemctl start raspimidihub-hostapd`. We never run at boot
# independently — raspimidihub decides when the AP comes up.
