# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Reload udev rules
        udevadm control --reload-rules 2>/dev/null || true

        # Enable and start the service
        systemctl daemon-reload
        systemctl enable raspimidihub.service
        systemctl start raspimidihub.service || true
        ;;
esac

#DEBHELPER#

exit 0
