# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Reload udev rules
        udevadm control --reload-rules 2>/dev/null || true

        # Unmask hostapd (Debian masks it by default) and disable +
        # stop hostapd / dnsmasq — our service manages these itself.
        # `--no-reload` per op + the single daemon-reload below
        # collapses the per-call "Synchronizing state of …" noise.
        systemctl --no-reload unmask hostapd 2>/dev/null || true
        systemctl --no-reload disable hostapd 2>/dev/null || true
        systemctl stop hostapd 2>/dev/null || true
        systemctl stop dnsmasq 2>/dev/null || true
        systemctl --no-reload disable dnsmasq 2>/dev/null || true

        # Set a UNIQUE hostname for mDNS discovery: raspimidihub-<wlan0
        # MAC suffix> (e.g. raspimidihub-735C) — the same token as the AP
        # SSID and the UI title bar, so a hub has one identity everywhere.
        # A bare "raspimidihub" collides between hubs on one network
        # (avahi renames the loser to raspimidihub-2.local). Keep this in
        # sync with wifi.default_hostname() / _get_mac_suffix().
        WLAN_MAC=$(cat /sys/class/net/wlan0/address 2>/dev/null | tr -d ':')
        if [ -n "$WLAN_MAC" ]; then
            SUFFIX=$(printf '%s' "$WLAN_MAC" | tail -c 4 | tr 'a-f' 'A-F')
            DESIRED="raspimidihub-$SUFFIX"
        else
            DESIRED="raspimidihub"
        fi
        if [ "$(hostname)" != "$DESIRED" ]; then
            hostnamectl set-hostname "$DESIRED" 2>/dev/null || true
            # Update /etc/hosts to match
            if grep -q "127.0.1.1" /etc/hosts; then
                sed -i "s/127\.0\.1\.1.*/127.0.1.1\t\t$DESIRED/" /etc/hosts
            else
                echo "127.0.1.1		$DESIRED" >> /etc/hosts
            fi
        fi

        # Create config directory on boot partition
        mkdir -p /boot/firmware/raspimidihub 2>/dev/null || true

        # Phase 5.5 update flow stores downloaded debs here so they
        # survive across reboots. Created here (not in the Makefile)
        # because it must exist on the live read-only system, not on
        # the build host.
        mkdir -p /var/lib/raspimidihub/updates 2>/dev/null || true

        # BlueZ's built-in `midi` plugin is disabled via a systemd
        # drop-in (etc/systemd/system/bluetooth.service.d/no-midi.conf
        # passes -P midi to bluetoothd). It creates a stub ALSA seq
        # client per BLE-MIDI device but doesn't forward GATT
        # notifications — verified with btmon + aseqdump on WIDI
        # Master. We bridge BLE-MIDI ourselves in
        # raspimidihub.ble_midi_bridge instead. (An earlier attempt
        # to disable via main.conf `DisablePlugins = midi` was a
        # no-op: BlueZ 5.82 prints "Unknown key DisablePlugins".)
        if [ -f /etc/bluetooth/main.conf ]; then
            sed -i '/^DisablePlugins.*midi/d' /etc/bluetooth/main.conf
        fi

        # BlueZ persistent state on read-only root: tmpfs at
        # /var/lib/bluetooth restored from a snapshot tarball on
        # /boot/firmware. Activated by the raspimidihub-bt-state
        # oneshot service ordered Before=bluetooth.service.
        systemctl --no-reload enable raspimidihub-bt-state.service \
            2>/dev/null || true

        # Restart bluetooth so the new drop-ins (no-midi + tmpfs
        # state mount) take effect right away rather than waiting
        # for the next reboot.
        systemctl restart raspimidihub-bt-state.service 2>/dev/null || true
        systemctl restart bluetooth.service 2>/dev/null || true

        # Enable BLE-MIDI profile in bluealsa (if installed). bluealsa
        # ships a2dp / hfp profiles by default; the `midi` profile is
        # what makes a connected BLE-MIDI device expose itself via ALSA
        # seq, which the engine then picks up via its normal scan.
        # Drop-in (not a file edit) so the package's stock unit stays
        # untouched and an apt purge cleanly reverts. Also unmask +
        # enable + start bluetooth.service and bluealsa.service —
        # an earlier system-prepare run may have disabled them as
        # part of the "MIDI-only appliance" trim.
        if [ -x /usr/bin/bluealsa ]; then
            mkdir -p /etc/systemd/system/bluealsa.service.d
            cat > /etc/systemd/system/bluealsa.service.d/midi.conf << 'CONF'
[Service]
ExecStart=
ExecStart=/usr/bin/bluealsa -S -p a2dp-source -p a2dp-sink -p midi
CONF
            systemctl --no-reload unmask bluetooth bluealsa 2>/dev/null || true
            systemctl daemon-reload
            systemctl enable --now bluetooth 2>/dev/null || true
            systemctl restart bluealsa 2>/dev/null || true
        fi

        # Trim system services + reserve a CPU core. Idempotent: safe
        # to run on every install / upgrade. Service / timer disable
        # takes effect immediately; CPU pin (AllowedCPUs=3) takes
        # effect on the next service restart below; kernel core
        # isolation (isolcpus=3) takes effect on next reboot.
        if [ -x /usr/local/bin/raspimidihub-system-prepare ]; then
            /usr/local/bin/raspimidihub-system-prepare || true
        fi

        # Enable and (re)start the service
        systemctl daemon-reload
        systemctl enable raspimidihub.service
        systemctl restart raspimidihub.service || true
        ;;
esac

#DEBHELPER#

exit 0
