# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Reload udev rules
        udevadm control --reload-rules 2>/dev/null || true

        # Unmask hostapd (Debian masks it by default) and stop default dnsmasq
        # Our service manages these itself
        systemctl unmask hostapd 2>/dev/null || true
        systemctl disable hostapd 2>/dev/null || true
        systemctl stop hostapd 2>/dev/null || true
        systemctl stop dnsmasq 2>/dev/null || true
        systemctl disable dnsmasq 2>/dev/null || true

        # Set hostname for mDNS discovery (raspimidihub.local)
        if [ "$(hostname)" != "raspimidihub" ]; then
            hostnamectl set-hostname raspimidihub 2>/dev/null || true
            # Update /etc/hosts to match
            if grep -q "127.0.1.1" /etc/hosts; then
                sed -i 's/127\.0\.1\.1.*/127.0.1.1\t\traspimidihub/' /etc/hosts
            else
                echo "127.0.1.1		raspimidihub" >> /etc/hosts
            fi
        fi

        # Create config directory on boot partition
        mkdir -p /boot/firmware/raspimidihub 2>/dev/null || true

        # Enable and (re)start the service
        systemctl daemon-reload
        systemctl enable raspimidihub.service
        systemctl restart raspimidihub.service || true
        ;;
esac

#DEBHELPER#

exit 0
